# deploy/macos/com.fallow.agent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  LaunchAgent (NOT a LaunchDaemon) for the Fallow per-machine agent.

  WHY A LaunchAgent, NOT A LaunchDaemon:
  Fallow's idle detection reads the console user's HID input timer
  (CGEventSourceSecondsSinceLastEventType via pyobjc Quartz). That API only
  returns meaningful values inside a logged-in GUI (Aqua) session. A
  LaunchDaemon runs in the system context (session 0) with no window server
  connection, so idle detection would always read "idle" and Fallow would never
  yield to the user. Therefore this MUST be a per-user LaunchAgent loaded into
  the user's GUI session (gui/$UID) via `launchctl bootstrap`.

  It is loaded into the console user's own GUI domain (gui/$UID) by install.sh,
  so the agent runs AS the pilot account — same uid, same login session, no
  privilege escalation. There is no UserName key: a LaunchAgent inherits the
  identity of the domain it is bootstrapped into.

  install.sh templates the __TOKENS__ below and writes the result to
  ~/Library/LaunchAgents/com.fallow.agent.plist. Do not load this file directly.
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.fallow.agent</string>

    <key>ProgramArguments</key>
    <array>
        <string>__PYTHON__</string>
        <string>-m</string>
        <string>fallow_agent</string>
        <string>run</string>
        <string>--config</string>
        <string>__CONFIG__</string>
    </array>

    <!-- Restart the agent whenever it exits (crash, preemption escalation,
         manual kill). Combined with RunAtLoad this gives boot + login start. -->
    <key>KeepAlive</key>
    <true/>
    <key>RunAtLoad</key>
    <true/>

    <!-- Backend selection lives here. install.sh detects the CPU arch, points
         FALLOW_LLAMA_SERVER_BINARY at the matching llama-server build (the Metal
         build on Apple Silicon, a CPU build on Intel), and on the CPU path sets
         a conservative OMP_NUM_THREADS cap so the agent stays gentle on a shared
         machine. The child llama-server inherits this environment. -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>FALLOW_LLAMA_SERVER_BINARY</key>
        <string>__LLAMA_BINARY__</string>
        <key>OMP_NUM_THREADS</key>
        <string>__OMP_THREADS__</string>
    </dict>

    <!-- Background priority: the agent must never contend with the user. -->
    <key>ProcessType</key>
    <string>Background</string>

    <!-- Throttle restart storms if the process exits immediately. -->
    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>StandardOutPath</key>
    <string>__STDOUT__</string>
    <key>StandardErrorPath</key>
    <string>__STDERR__</string>

    <key>WorkingDirectory</key>
    <string>__WORKDIR__</string>
</dict>
</plist>
